# autobot-celery.conf
[program:autobot-celery-flower]
command=/home/leo/Env/autobot/bin/celery flower -A autobot --loglevel=INFO
directory=/home/leo/autobot
user=leo
numprocs=1
stdout_logfile=/home/leo/Env/autobot/logs/celery_flower.log
stderr_logfile=/home/leo/Env/autobot/logs/celery_flower.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000

[program:autobot-celery-beat]
command=/home/leo/Env/autobot/bin/celery -A autobot worker -l info -B -n beat_worker@ubuntu --purge
directory=/home/leo/autobot
user=leo
numprocs=1
stdout_logfile=/home/leo/Env/autobot/logs/celery_beat.log
stderr_logfile=/home/leo/Env/autobot/logs/celery_beat.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000

[program:autobot-celery-ad-bot-runner]
command=/home/leo/Env/autobot/bin/celery -A autobot -Q ad_bot_runner worker -l info -n ad_bot_runner@ubuntu --purge
directory=/home/leo/autobot
user=leo
numprocs=1
stdout_logfile=/home/leo/Env/autobot/logs/celery_ad_bot_runner.log
stderr_logfile=/home/leo/Env/autobot/logs/celery_ad_bot_runner.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000

[program:autobot-celery-sell-ad-bot-execution]
command=/home/leo/Env/autobot/bin/celery -A autobot -Q sell_ad_bot_execution worker -l info -n sell_ad_bot_execution@ubuntu --purge
directory=/home/leo/autobot
user=leo
numprocs=1
stdout_logfile=/home/leo/Env/autobot/logs/celery_sell_ad_bot_execution.log
stderr_logfile=/home/leo/Env/autobot/logs/celery_sell_ad_bot_execution.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000

[program:autobot-celery-buy-ad-bot-execution]
command=/home/leo/Env/autobot/bin/celery -A autobot -Q buy_ad_bot_execution worker -l info -n buy_ad_bot_execution@ubuntu --purge
directory=/home/leo/autobot
user=leo
numprocs=1
stdout_logfile=/home/leo/Env/autobot/logs/celery_buy_ad_bot_execution.log
stderr_logfile=/home/leo/Env/autobot/logs/celery_buy_ad_bot_execution.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000

[program:autobot-celery-seller-bot-handler]
command=/home/leo/Env/autobot/bin/celery -A autobot -Q seller_bot_handler worker -l info -n seller_bot_handler@ubuntu --purge
directory=/home/leo/autobot
user=leo
numprocs=1
stdout_logfile=/home/leo/Env/autobot/logs/celery_seller_bot_handler.log
stderr_logfile=/home/leo/Env/autobot/logs/celery_seller_bot_handler.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000

[program:autobot-celery-fast-rare-tasks]
command=/home/leo/Env/autobot/bin/celery -A autobot -Q fast_rare_tasks worker -l info -n fast_rare_tasks@ubuntu --purge
directory=/home/leo/autobot
user=leo
numprocs=1
stdout_logfile=/home/leo/Env/autobot/logs/celery_fast_rare_tasks.log
stderr_logfile=/home/leo/Env/autobot/logs/celery_fast_rare_tasks.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000

[program:autobot-celery-qiwi-status-updater]
command=/home/leo/Env/autobot/bin/celery -A autobot -Q qiwi_status_updater worker -l info -n qiwi_status_updater@ubuntu --purge
directory=/home/leo/autobot
user=leo
numprocs=1
stdout_logfile=/home/leo/Env/autobot/logs/celery_qiwi_status_updater.log
stderr_logfile=/home/leo/Env/autobot/logs/celery_qiwi_status_updater.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000

[program:autobot-celery-qiwi-profit-fixator]
command=/home/leo/Env/autobot/bin/celery -A autobot -Q qiwi_profit_fixator worker -l info -n qiwi_profit_fixator@ubuntu --purge
directory=/home/leo/autobot
user=leo
numprocs=1
stdout_logfile=/home/leo/Env/autobot/logs/celery_qiwi_profit_fixator.log
stderr_logfile=/home/leo/Env/autobot/logs/celery_qiwi_profit_fixator.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000

[group:autobot-celery-groups]
programs=autobot-celery-flower,autobot-celery-beat,autobot-celery-ad-bot-runner,autobot-celery-sell-ad-bot-execution,autobot-celery-buy-ad-bot-execution,autobot-celery-seller-bot-handler,autobot-celery-fast-rare-tasks,autobot-celery-qiwi-status-updater,autobot-celery-qiwi-profit-fixator
priority=1000
